# recipes-atlas-soc/atlas-soc-fftsw-apps-init/atlas-soc-fftsw-apps-init/atlas-soc-fftsw-init.service
[Unit]
Description=Atlas SoC fftsw init

[Service]
ExecStart=/examples/fft/bin/setup_target_fft_env.sh

[Install]
WantedBy=multi-user.target
